# apache2.service
[Service]
ExecStart=/usr/sbin/apache2 -k start